# files/etc/profile.d/homebrew.sh
#!/bin/sh

# Homebrew
export HOMEBREW_HOME="/home/linuxbrew/.linuxbrew"
export HOMEBREW_CACHE="${CACHE_HOME}/homebrew"
export PATH="${HOMEBREW_HOME}/bin:${HOMEBREW_HOME}/sbin:${PATH}"
