# files/etc/profile.d/packer.sh
#!/bin/sh

# Packer Environment Variables: https://www.packer.io/docs/other/environment-variables.html
export CHECKPOINT_DISABLE=1
